# packages/lcdproc/main/debian/jammy/lcdproc.service
[Unit]
Description=LCD display daemon
Documentation=man:LCDd(8) http://www.lcdproc.org/

[Service]
User=root
ExecStart=/usr/sbin/LCDd -s 1 -f -c /etc/LCDd.conf

[Install]
WantedBy=multi-user.target